# services/xstrm.service
[Unit]
Description=XSTRM interactive CLI helper
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/xstrm-suite
ExecStart=/opt/xstrm-suite/bin/xstrm
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
